# src/crontest.cron
00 11 * * * /home/zkr/pyWorkSpace/For_Majesty/reconstruct/src/main.sh
